# automation/sms-data-feed.service
[Unit]
Description=SMS Spam Detection — Real-time Data Feed (continuous ingestion + pipeline trigger)
After=network.target

[Service]
Type=simple
User=azureuser
WorkingDirectory=/home/azureuser/SMS-Spam-Detection
ExecStart=/home/azureuser/SMS-Spam-Detection/.venv/bin/python automation/data_feed.py
Restart=always
RestartSec=10
StandardOutput=append:/home/azureuser/SMS-Spam-Detection/logs/data_feed.log
StandardError=append:/home/azureuser/SMS-Spam-Detection/logs/data_feed.log
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
